# nginx-passenger/asset/supervisord.conf
[supervisord]
nodaemon = true
user     = root

[program:nginx]
command=/opt/nginx/sbin/nginx
stdout_events_enabled=true
stderr_events_enabled=true