# ampache/asset/supervisord.conf
[supervisord]
nodaemon = true
user     = root

[program:apache2]
command         = /usr/sbin/apache2ctl -D "FOREGROUND" -k start
environment     = APACHE_LOG_DIR=/var/log/apache2
redirect_stderr = true