# postfix-gmail/asset/supervisord.conf
[supervisord]
nodaemon = true
user     = root

[program:postfix]
command	= /opt/postfix.sh
